# 08-scheduling/macos/compound-review.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" 
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.yourproject.daily-compound-review</string>
  
  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>/Users/you/projects/your-project/scripts/daily-compound-review.sh</string>
  </array>
  
  <key>WorkingDirectory</key>
  <string>/Users/you/projects/your-project</string>
  
  <key>StartCalendarInterval</key>
  <dict>
    <key>Hour</key>
    <integer>22</integer>
    <key>Minute</key>
    <integer>30</integer>
  </dict>
  
  <key>StandardOutPath</key>
  <string>/Users/you/projects/your-project/logs/compound-review.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/you/projects/your-project/logs/compound-review.log</string>
  
  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
  </dict>
</dict>
</plist>
